# launchagents/com.kalyan.gbh.kovacs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
    <key>Label</key><string>com.kalyan.gbh.kovacs</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3.11</string>
        <string>/Users/kalyan/Documents/projects/gbh/main.py</string>
        <string>kovacs</string>
    </array>
    <key>WorkingDirectory</key><string>/Users/kalyan/Documents/projects/gbh</string>
    <key>StartCalendarInterval</key>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>0</integer></dict>
    <key>StandardOutPath</key><string>/tmp/gbh_kovacs.out</string>
    <key>StandardErrorPath</key><string>/tmp/gbh_kovacs.err</string>
</dict></plist>
